# server/etc/services.d/netxmsd/run
#!/usr/bin/execlineb -P
/usr/bin/netxmsd -c /opt/netxmsd/netxmsd.conf -q